# deploy/systemd/amf-conversation-memory-extractor.service
[Unit]
Description=Agent Memory Fabric conversation-memory extractor (quality-gated template)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=amf-extractor
Group=amf-extractor
WorkingDirectory=/opt/agent-memory-fabric
ExecCondition=/usr/bin/node /opt/agent-memory-fabric/scripts/amf-conversation-memory-extractor.mjs --config /etc/agent-memory-fabric/conversation-memory-extractor.json --verify-quality-gate
ExecStart=/usr/bin/flock -n /var/lib/agent-memory-fabric/conversation-memory-extractor/tick.lock /usr/bin/node /opt/agent-memory-fabric/scripts/amf-conversation-memory-extractor.mjs --config /etc/agent-memory-fabric/conversation-memory-extractor.json --scheduled
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/agent-memory-fabric/conversation-memory-extractor
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictSUIDSGID=true
