# deploy/systemd/amf-raw-memory-extractor.timer
[Unit]
Description=Run slow Agent Memory Fabric RAW-to-memory extractor

[Timer]
OnBootSec=4min
OnUnitActiveSec=5min
RandomizedDelaySec=20s
AccuracySec=5s
Unit=amf-raw-memory-extractor.service

[Install]
WantedBy=timers.target
